# templates/supervisor.conf
[program:{name}]
command=fig up
directory=/web/apps/{name}
